# wg-resrolve-dns.service
[Unit]
Description=WireGuard Tools - Reresolve DNS

[Service]
Type=oneshot
ExecStart=/bin/bash /usr/local/bin/reresolve-dns.sh
