# wg-reresolve-dns.timer
[Unit]
Description=WireGuard Tools - Reresolve DNS

[Timer]
OnUnitActiveSec=60s
OnBootSec=300s

[Install]
WantedBy=timers.target
